# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=python app.py

[program:rails]
command=rails s -b 0.0.0.0