# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=/root/fast-style-transfer/python app.py

[program:rails]
command=/root/mod-image-server/rails s -b 0.0.0.0